# files/student/etc/systemd/system/reset-opilas-home.service
# Resets /home/opilas to /etc/skel. Started (blocking) by GDM PostSession on
# every logout; can also be run by hand: `systemctl start reset-opilas-home`.
# Enabled in build.student.sh mainly so it's present and startable — the actual
# trigger is PostSession, not a target dependency.
[Unit]
Description=Reset /home/opilas to a pristine /etc/skel
# Runs at boot too, to create the home before opilas autologs in (bootc doesn't
# persist a build-time /home). Ordered after /var is mounted, after sysusers has
# created the opilas account (declared in usr/lib/sysusers.d/svk-opilas.conf, so
# the chown target exists), and before the display manager.
After=local-fs.target var-home.mount systemd-sysusers.service
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/reset-opilas-home.sh

[Install]
# Enabled in build.student.sh -> runs every boot; PostSession also starts it on
# each logout.
WantedBy=multi-user.target
