# files/student/etc/systemd/system/reset-student-home.service
# Resets /home/student to /etc/skel. Started (blocking) by GDM PostSession on
# every logout; can also be run by hand: `systemctl start reset-student-home`.
# Enabled in build.student.sh mainly so it's present and startable — the actual
# trigger is PostSession, not a target dependency.
[Unit]
Description=Reset /home/student to a pristine /etc/skel
# Runs at boot too, to create the home before the student autologs in (bootc
# doesn't persist a build-time /home). Ordered after /var is mounted and before
# the display manager so the home exists by the time GDM logs the student in.
After=local-fs.target var-home.mount
Before=display-manager.service

[Service]
Type=oneshot
ExecStart=/usr/bin/reset-student-home.sh

[Install]
# Enabled in build.student.sh -> runs every boot; PostSession also starts it on
# each logout.
WantedBy=multi-user.target
