# files/base/etc/systemd/system/bootc-fetch-apply-updates.service.d/10-svk-no-force-apply.conf
# Upstream's unit runs `bootc upgrade --apply --quiet`, which reboots the
# machine the instant an update is staged. At svk's monthly/weekly per-flavor
# cadence (see files/{student,staff,server}/…/bootc-fetch-apply-updates.timer.d/)
# that could land a forced reboot mid-class or mid-workday whenever a boot
# happens to coincide with the scheduled window — exactly the kind of
# unpredictable disruption the whole cadence design is meant to avoid. Replace
# ExecStart with our own wrapper: stages only (no --apply, applies at the
# machine's own next natural boot) and logs the run to /var/log/svk/update.log.
#
# Desktop-only override — svk-server's own timer.d drop-in does NOT set
# --apply=false; the server keeps upstream's force-apply behaviour on purpose
# (see tasks/todo/…-fleet-update-cadence.md).
[Service]
ExecStart=
ExecStart=/usr/libexec/svk/svk-update-image
