# files/student/etc/systemd/system/bootc-fetch-apply-updates.timer.d/10-svk-schedule.conf
# svk-student cadence: monthly, first Monday of the month, same time-of-day as
# staff/server (one consistent "update morning" flavor across the fleet even
# though frequency differs). Locked-down kiosk (autologin, home reset every
# logout, curated/restricted app set, no free browsing) is genuinely lower risk
# than a staff desktop, so a slower patch cadence is an acceptable trade for far
# fewer disruptive windows on machines a class is actively using. See
# tasks/todo/…-fleet-update-cadence.md.
#
# OnBootSec=/OnUnitInactiveSec= are cumulative directives in systemd — an
# empty assignment clears upstream's OnBootSec=1h/OnUnitInactiveSec=8h before
# OnCalendar= below takes over as the sole trigger. RandomizedDelaySec/
# Persistent are single-value and simply override upstream's.
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=Mon *-*-1..7 05:00:00
RandomizedDelaySec=1d
Persistent=true
