# archived/files/base/etc/systemd/system/bootc-fetch-apply-updates.timer.d/10-svk-randomize.conf
# Anti-thundering-herd jitter for the bootc updater, in case this base drives
# updates through bootc's own timer rather than uupd. Harmless if the unit isn't
# active (systemd just has nothing to apply it to). See also the uupd.timer
# drop-in and TODO §11.
[Timer]
RandomizedDelaySec=2h
Persistent=true
