# files/staff/etc/systemd/system/bootc-fetch-apply-updates.timer.d/10-svk-schedule.conf
# svk-staff cadence: weekly, Monday morning before the school day starts.
# Staff have real desktop freedom (browsing, their own workflow) unlike the
# locked-down student kiosk, so a tighter cadence than student's monthly is the
# right risk trade. See tasks/todo/…-fleet-update-cadence.md.
#
# OnBootSec=/OnUnitInactiveSec= are cumulative directives in systemd — an
# empty assignment clears upstream's OnBootSec=1h/OnUnitInactiveSec=8h before
# OnCalendar= below takes over as the sole trigger. RandomizedDelaySec/
# Persistent are single-value and simply override upstream's.
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=Mon *-*-* 05:00:00
RandomizedDelaySec=3h
Persistent=true
