# files/server/etc/systemd/system/bootc-fetch-apply-updates.timer.d/10-svk-schedule.conf
# svk-server cadence: same weekly Monday-morning window as the desktops — one
# fleet-wide maintenance window is easier to reason about than a second
# schedule. Unlike the desktops, the server keeps upstream's --apply behaviour:
# no service.d ExecStart= override here, so build.server.sh's
# `systemctl enable bootc-fetch-apply-updates.timer` still runs the unmodified
# `bootc upgrade --apply --quiet`. A server is more likely to go a long time
# between organic reboots, so stage-and-wait could leave it un-patched far
# longer than intended, and a brief scheduled blip to mirror/DNS/cache at 5am
# Monday is an acceptable, predictable cost since nothing else in the fleet is
# active then either. See tasks/todo/…-fleet-update-cadence.md.
#
# OnBootSec=/OnUnitInactiveSec= are cumulative directives in systemd — an
# empty assignment clears upstream's OnBootSec=1h/OnUnitInactiveSec=8h before
# OnCalendar= below takes over as the sole trigger. RandomizedDelaySec/
# Persistent are single-value and simply override upstream's.
[Timer]
OnBootSec=
OnUnitInactiveSec=
OnCalendar=Mon *-*-* 05:00:00
RandomizedDelaySec=3h
Persistent=true
