# archived/files/base/etc/systemd/system/uupd.timer.d/10-svk-randomize.conf
# Spread desktop auto-updates across the fleet so ~19 laptops don't all hit the
# cache (and its upstream ghcr.io pull) at the same instant. Bluefin's updater is
# uupd; this adds a large random jitter to whenever its timer fires, without
# changing the base schedule.
#
# If your Bluefin actually uses a different updater (check `systemctl list-timers`
# — it may be bootc-fetch-apply-updates.timer or rpm-ostreed-automatic.timer),
# move this jitter to that unit's drop-in. See TODO §11.
[Timer]
RandomizedDelaySec=2h
