# archived/files/server/etc/systemd/system/bootc-fetch-apply-updates.timer.d/10-svk-randomize.conf
# The server (uCore) updates via bootc's own timer, which build.server.sh
# enables. Give it a daily schedule with a big random jitter — there's only one
# server so the jitter matters less, but keep it consistent with the desktops.
[Timer]
OnCalendar=
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=2h
Persistent=true
